# services/memcached/run
#!/bin/bash

exec /usr/bin/memcached -u abc -v 